# config-files/update_bootinfo.service.in
[Unit]
Description=Record successful system boot
Requires=boot-complete.target
After=boot-complete.target
ConditionPathExists=!@CMAKE_INSTALL_SYSCONFDIR@/initrd-release

[Service]
Type=oneshot
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/tegra-bootinfo -b
RemainAfterExit=yes

[Install]
WantedBy=basic.target
