# config-files/update_bootinfo.service.in
[Unit]
Description=Record successful system boot
Requires=boot-complete.target
After=boot-complete.target
ConditionPathExists=!@sysconfdir@/initrd-release

[Service]
Type=oneshot
ExecStart=@bindir@/tegra-bootinfo -b
RemainAfterExit=yes

[Install]
WantedBy=basic.target
